# deploy/cai-web.service
[Unit]
Description=CAI — cai.canine.dev (the open standard: rubric API + UI + scorer)
Documentation=https://github.com/CanineCC/CAI
After=network-online.target
Wants=network-online.target

# Type=notify + WatchdogSec give the service a readiness signal and a liveness probe: the app (UseSystemd) sends
# sd_notify READY=1 once it is actually serving, so a restart/deploy is only "up" when requests can be handled, and
# emits WATCHDOG=1 keepalives — systemd restarts the process if they stop (a hung-but-alive app is caught, not just a
# crashed one). Restart=on-failure + RestartSec make that automatic and reversible.
[Service]
Type=notify
NotifyAccess=main
User=jimmy
Group=jimmy
WorkingDirectory=/home/jimmy/apps/cai-web/app
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://0.0.0.0:8090
Environment=HOME=/home/jimmy
Environment=Rubrics__Root=/home/jimmy/apps/cai-web/rubrics
ExecStart=/usr/bin/dotnet /home/jimmy/apps/cai-web/app/Cai.Web.dll
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=0
# Liveness: the app must check in within this window or systemd restarts it.
WatchdogSec=30s

[Install]
WantedBy=multi-user.target
